# systemd/.systemd-root/update-seterc.service
[Unit]
Description=Update seterc values

[Service]
Type=oneshot
ExecStart=/opt/bin/update-seterc

[Install]
WantedBy=multi-user.target
