# systemd/.config/systemd/user/xhack.target
[Unit]
Description=A Hack Target to start user apps only after X starts.
